# files/shared/etc/profile.d/umbra-aliases.sh
#!/usr/bin/env bash

# All the aliases

alias fastfetch="/usr/bin/fastfetch --config /usr/share/umbra/configs/fastfetch/config.jsonc"
