# files/shared/etc/xdg/autostart/umbra-theme-manager.desktop
[Desktop Entry]
Type=Application
Name=Nova Theme Manager
Exec=/usr/bin/systemctl --user enable umbra-theme-manager
