# files/shared/etc/xdg/autostart/umbra-theme-manager.desktop
[Desktop Entry]
Comment[en_US]=
Comment=
Exec=/usr/bin/systemctl --user enable umbra-theme-manager
